# configs/debian/cron.d/imscp
# Quota accounting
0,30 * * * * root {QUOTA_ROOT_DIR}/imscp-dsk-quota &>{LOG_DIR}/imscp-dsk-quota.log

# Traffic accounting
0,30 * * * * root {TRAFF_ROOT_DIR}/imscp-srv-traff &>{LOG_DIR}/imscp-srv-traff.log
0,30 * * * * root {TRAFF_ROOT_DIR}/imscp-vrl-traff &>{LOG_DIR}/imscp-vrl-traff.log

# Customer logs
@daily root {TOOLS_ROOT_DIR}/imscp-httpd-logs-mngr &>{LOG_DIR}/imscp-httpd-logs-mngr.log

# Backup customers' data depending of the domain properties
{BACKUP_MINUTE} {BACKUP_HOUR} * * * root {CMD_NICE} -n 10 {BACKUP_ROOT_DIR}/imscp-backup-all &>{LOG_DIR}/imscp-backup-all-mngr.log

# Backup i-MSCP configuration files and database
@daily root {CMD_NICE} -n 10 {BACKUP_ROOT_DIR}/imscp-backup-imscp &>{LOG_DIR}/imscp-backup-imscp-mngr.log

# Remove configuration backups older than 7 days
@weekly root {CMD_FIND} {CONF_DIR}/*/backup -type -f -mtime +7 -regextype sed -regex '.*/.*[0-9]\{10\}$' -exec {CMD_RM} -- {} +

# Remove backend Logs older than 7 days
@weekly root {CMD_FIND} {LOG_DIR} -type f -mtime +7 -exec {CMD_RM} -- {} +

# Disable expired accounts
@daily root {TOOLS_ROOT_DIR}/imscp-disable-accounts >/dev/null 2>&1

# imscp [{ENTRY_ID}] entry BEGIN
# imscp [{ENTRY_ID}] entry ENDING
