# configs/debian/cron.d/imscp

# delayed tasks START.

# Quota
@daily root umask 027; {QUOTA_ROOT_DIR}/imscp-dsk-quota &>{LOG_DIR}/imscp-dsk-quota.log

# Traffic
0,30	*	*	*	*  root umask 027; {TRAFF_ROOT_DIR}/imscp-srv-traff &>{LOG_DIR}/imscp-srv-traff.log
0,30	*	*	*	*  root umask 027; {TRAFF_ROOT_DIR}/imscp-vrl-traff &>{LOG_DIR}/imscp-vrl-traff.log

# customer logs
@daily root umask 027; {TOOLS_ROOT_DIR}/imscp-httpd-logs-mngr &>{LOG_DIR}/imscp-httpd-logs-mngr.log

# Backup for all the customers' data depending of the domain properties
@daily root umask 027; {BACKUP_ROOT_DIR}/imscp-backup-all yes &>{LOG_DIR}/imscp-backup-all-mngr.log

# Backup for i-MSCP database and all /etc/imscp/* configuration files
@daily root umask 027; {BACKUP_ROOT_DIR}/imscp-backup-imscp noreport &>{LOG_DIR}/imscp-backup-imscp-mngr.log

# Remove config backups older than seven days
@daily root find {CONF_DIR}/*/backup/* -maxdepth 0 -type f -mtime +7 -print | egrep '.*\.[0-9]+$' | xargs -r /bin/rm

# Remove Daemon Logs older than 14 days (except .gz files)
@daily root find {LOG_DIR}/* -maxdepth 1 -type f -mtime +14 -print | egrep '.*\.gz$' | xargs -r /bin/rm

# AWStats
{AW-ENABLED}15	*/6	*	*	*	root {AWSTATS_ROOT_DIR}/awstats_updateall.pl now -awstatsprog={AWSTATS_ENGINE_DIR}/awstats.pl >/dev/null 2>&1

# Rootkit Hunter
{RK-ENABLED}@weekly	root {RKHUNTER} --cronjob --createlogfile {RKHUNTER_LOG}.root --display-logfile 1>{RKHUNTER_LOG} 2>/dev/null

# Chkrootkit
{CR-ENABLED}@weekly	root {CHKROOTKIT} 1>{CHKROOTKIT_LOG} 2>&1

# Look for and purge old sessions every 30 minutes
0,30	*	*	*	*	root	{TOOLS_ROOT_DIR}/imscp-del-tmp >/dev/null 2>&1

# Delete old apache logfiles weekly
@weekly root {TOOLS_ROOT_DIR}/imscp-del-logs >/dev/null 2>&1

# [{DMN_NAME}]:{CRONJOB_ID} custom task START.
# [{DMN_NAME}]:{CRONJOB_ID} custom task END.

# [{TASKID}] custom task START.
# [{TASKID}] custom task END.

# delayed tasks END.
