# configs/centos/init.d/ispcp_daemon
#!/bin/sh

# ispCP ω (OMEGA) a Virtual Hosting Control System
# Copyright (c) 2007 by ispCP
# http://www.isp-control.net
#
#
# License:
#    This program is free software; you can redistribute it and/or
#    modify it under the terms of the MPL Mozilla Public License
#    as published by the Free Software Foundation; either version 1.1
#    of the License, or (at your option) any later version.
#
#    This program is distributed in the hope that it will be useful,
#    but WITHOUT ANY WARRANTY; without even the implied warranty of
#    MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#    MPL Mozilla Public License for more details.
#
#    You may have received a copy of the MPL Mozilla Public License
#    along with this program.
#
#    An on-line copy of the MPL Mozilla Public License can be found
#    http://www.mozilla.org/MPL/MPL-1.1.html
#
#
# The ispCP Home Page is at:
#
#    http://www.isp-control.net
#
### BEGIN INIT INFO
# Provides:             ispcp_daemon
# Required-Start:       $network $local_fs $remote_fs
# Required-Stop:
# Should-Stop:          $local_fs
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    ispCP GUI-Backend communication Daemon
### END INIT INFO
# Note: do not modify any of these vars here, use /etc/default/$NAME instead



# Source function library.
if [ -f /etc/init.d/functions ] ; then
  . /etc/init.d/functions
elif [ -f /etc/rc.d/init.d/functions ] ; then
  . /etc/rc.d/init.d/functions
else
  exit 0
fi

# Avoid using root's TMPDIR
unset TMPDIR

# Source networking configuration.
. /etc/sysconfig/network

if [ -f /var/www/ispcp/daemon ]; then
   . /var/www/ispcp/daemon
fi

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

# Check that ispcp.conf exists.
[ -f /etc/ispcp/ispcp.conf ] || exit 0

RETVAL=0


start() {
        KIND="ISPCP"
	echo -n $"Starting $KIND services: "
	daemon /var/www/ispcp/daemon/ispcp_daemon
	RETVAL=$?
	echo
      	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/ispcp || RETVAL=1
	echo ""
	return $RETVAL
}	

stop() {
        KIND="ISPCP"
	echo -n $"Shutting down $KIND services: "
	killproc ispcp_daemon
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/ispcp
	echo ""
	return $RETVAL
}	

restart() {
	stop
	start
}	

reload() {
        echo -n $"Reloading ispcp.conf: "
	killproc ispcp_daemon -HUP
	RETVAL=$?
	echo
	return $RETVAL
}	

rhstatus() {
	status ispcp_daemon
	RETVAL=$?
	if [ $? -ne 0 -o $RETVAL -ne 0 ] ; then
		return 1
	fi
}	


# Allow status as non-root.
if [ "$1" = status ]; then
       rhstatus
       exit $?
fi

# Check that we can write to it... so non-root users stop here
[ -w /etc/ispcp/ispcp.conf ] || exit 0



case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  restart)
  	restart
	;;
  reload)
  	reload
	;;
  status)
  	rhstatus
	;;
  condrestart)
  	[ -f /var/lock/subsys/ispcp ] && restart || :
	;;
  *)
	echo $"Usage: $0 {start|stop|restart|reload|status|condrestart}"
	exit 1
esac

exit $?
