# configs/debian/default/systemd/imscp_traffic.service
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
Description=i-MSCP traffic logger
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/perl /var/www/imscp/engine/tools/imscp-net-traffic-logger start
ExecStop=/usr/bin/perl /var/www/imscp/engine/tools/imscp-net-traffic-logger stop

[Install]
WantedBy=multi-user.target
