# configs/debian/init.d/imscp_network
#!/bin/sh
### BEGIN INIT INFO
# Provides:          imscp_network
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs $network
# Should-Start:      mysql
# Should-Stop:       mysql
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: i-MSCP network manager
### END INIT INFO

# Note: do not modify any of these vars here, use /etc/default/$NAME instead

PATH=/sbin:/bin:/usr/sbin:/usr/bin
DESC="i-MSCP network manager"
NAME=imscp_network
LFILE=/var/run/$NAME
TOOLSDIR="/var/www/imscp/engine/tools"

. /lib/lsb/init-functions

# Read config file if present.
if [ -r /etc/default/$NAME ]; then
	. /etc/default/$NAME
fi

# Don't run if we are running upstart
if init_is_upstart 2>/dev/null; then
	exit 1
fi

do_start()
{
	${TOOLSDIR}/imscp-net-interfaces-mngr start 2>/dev/null
	${TOOLSDIR}/imscp-net-traffic-logger start 2>/dev/null

	# Touch lock file
	touch $LFILE
}

do_stop()
{
	${TOOLSDIR}/imscp-net-traffic-logger stop 2>/dev/null
	${TOOLSDIR}/imscp-net-interfaces-mngr stop 2>/dev/null

	# Remove lock file
	rm -f $LFILE
}

case "$1" in
  start)
	log_daemon_msg "Starting $DESC" "$NAME"

	if [ -e "$LFILE" ]; then
		log_progress_msg "already started"
		log_end_msg 0
	else
		do_start
		log_end_msg $?
	fi
	;;
  stop)
	log_daemon_msg "Stopping $DESC" "$NAME"
	if [ ! -e "$LFILE" ]; then
		log_progress_msg "already stopped"
		log_end_msg 0
	else
		do_stop
		log_end_msg $?
	fi
	;;
  restart|force-reload)
	log_daemon_msg "Stopping $DESC" "$NAME"

	if [ ! -e "$LFILE" ]; then
		log_progress_msg "already stopped"
		log_end_msg 0
	else
		do_stop
		log_end_msg $?
	fi

	log_daemon_msg "Starting $DESC" "$NAME"
	do_start
	log_end_msg $?
	;;
  status)
	if [ -e "$LFILE" ]; then
		log_success_msg "$NAME is running"
	else
		log_failure_msg "$NAME is not running"
		exit 3
	fi
	;;
  *)
	log_action_msg "Usage: /etc/init.d/$NAME {start|stop|restart|force-reload|status}"
	exit 1
	;;
esac
