# configs/freebsd/init.d/ispcp_daemon
#!/bin/sh
#
# $FreeBSD: ports/sysutils/ispcp/files/ispcp_daemon.sh.in,v 1.3 2010/02/13 20:47:00 fxtw Exp $
#

# PROVIDE: ispcp_daemon
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable ispcp_daemon:
#
#ispcp_daemon_enable="YES"
#

. /etc/rc.subr

name=ispcp_daemon
rcvar=`set_rcvar`

command=/usr/local/www/ispcp/daemon/ispcp_daemon
pidfile=/var/run/ispcp-daemon.pid

stop_postcmd=stop_postcmd

stop_postcmd() {
  rm -f $pidfile
}

# set defaults

ispcp_daemon_enable=${ispcp_daemon_enable:-"NO"}
ispcp_daemon_flags=${ispcp_daemon_flags:-""}

extra_commands="reload"

load_rc_config $name
run_rc_command "$1"