# configs/freebsd/init.d/ispcp_network
#!/usr/local/bin/bash

echo "Starting ispcp_network: currently not ported to BSD"
exit 0

# ispCP ω (OMEGA) a Virtual Hosting Control System
# Copyright (c) 2007 by ispCP
# http://www.isp-control.net
#
#
# License:
#    This program is free software; you can redistribute it and/or
#    modify it under the terms of the MPL Mozilla Public License
#    as published by the Free Software Foundation; either version 1.1
#    of the License, or (at your option) any later version.
#
#    This program is distributed in the hope that it will be useful,
#    but WITHOUT ANY WARRANTY; without even the implied warranty of
#    MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#    MPL Mozilla Public License for more details.
#
#    You may have received a copy of the MPL Mozilla Public License
#    along with this program.
#
#    An on-line copy of the MPL Mozilla Public License can be found
#    http://www.mozilla.org/MPL/MPL-1.1.html
#
#
# The ispCP Home Page is at:
#
#    http://www.isp-control.net
### BEGIN INIT INFO
# Provides:             ispcp_network
# Required-Start:       $network $local_fs $remote_fs
# Required-Stop:
# Should-Stop:          $local_fs
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    ispCP Network Traffic Logger
### END INIT INFO
# Note: do not modify any of these vars here, use /etc/default/$NAME instead

RMMOD="/sbin/rmmod"

IPTABLES="/sbin/iptables"
LOGDIR=/var/log/ispcp
LOGFILE=${LOGDIR}/${NAME}.log

${RMMOD} ipchains 1>/tmp/ispcp_network.log 2>/tmp/ispcp_network.err.log

# To monitor more ports, edit SERVICES variable add your own ports
# (ftp, proxy, http, etc.)
SERVICES=(
	# Web traffic
		80 	443
	# Pop 3
		110
	# IMAP
		143
	# Mail
		25 587 995 993
)

#
# firewall commet; if u have some kind of firewall
# then this hurts...
#
# must be chnaged in future releases!
# 1zo
#
# ${IPTABLES} -F INPUT
# ${IPTABLES} -F OUTPUT
# ${IPTABLES} -F FORWARD
#

# all traffic

${IPTABLES} -I INPUT

${IPTABLES} -I OUTPUT

# Services from matrix
for PORT in ${SERVICES[@]}; do
	${IPTABLES} -I ISPCP_INPUT  -p tcp --dport "$PORT" 2>> "$LOGFILE"
	${IPTABLES} -I ISPCP_OUTPUT -p tcp --sport "$PORT" 2>> "$LOGFILE"
done
