# configs/freebsd/init.d/ispcp_network
#!/usr/local/bin/bash

echo "Starting ispcp_network: currently not ported to BSD"
exit 0

# ispCP ω (OMEGA) a Virtual Hosting Control System
# Copyright (c) 2007 by ispCP
# http://www.isp-control.net
#
#
# License:
#    This program is free software; you can redistribute it and/or
#    modify it under the terms of the MPL Mozilla Public License
#    as published by the Free Software Foundation; either version 1.1
#    of the License, or (at your option) any later version.
#
#    This program is distributed in the hope that it will be useful,
#    but WITHOUT ANY WARRANTY; without even the implied warranty of
#    MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#    MPL Mozilla Public License for more details.
#
#    You may have received a copy of the MPL Mozilla Public License
#    along with this program.
#
#    An on-line copy of the MPL Mozilla Public License can be found
#    http://www.mozilla.org/MPL/MPL-1.1.html
#
#
# The ispCP Home Page is at:
#
#    http://www.isp-control.net
### BEGIN INIT INFO
# Provides:             ispcp_network
# Required-Start:       $network $local_fs $remote_fs
# Required-Stop:
# Should-Stop:          $local_fs
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    ispCP Network Traffic Logger
### END INIT INFO
# Note: do not modify any of these vars here, use /etc/default/$NAME instead

RMMOD='/sbin/rmmod'

IPTABLES="/sbin/iptables"

${RMMOD} ipchains 1>/tmp/ispcp_network.log 2>/tmp/ispcp_network.err.log

#
# firewall commet; if u have some kind of firewall
# then this hurts... 
#
# must be chnaged in future releases!
# 1zo
#
# ${IPTABLES} -F INPUT
# ${IPTABLES} -F OUTPUT
# ${IPTABLES} -F FORWARD
#

# all traffic

${IPTABLES} -I INPUT

${IPTABLES} -I OUTPUT

# web traffic

${IPTABLES} -I INPUT -p tcp --dport 80
${IPTABLES} -I OUTPUT -p tcp --sport 80

# pop3 traffic

${IPTABLES} -I INPUT -p tcp --dport 110
${IPTABLES} -I OUTPUT -p tcp --sport 110

# imap traffic

${IPTABLES} -I INPUT -p tcp --dport 143
${IPTABLES} -I OUTPUT -p tcp --sport 143

# mail traffic

${IPTABLES} -I INPUT  -p tcp --dport 25
${IPTABLES} -I OUTPUT -p tcp --sport 25

${IPTABLES} -I INPUT  -p tcp --dport 25
${IPTABLES} -I OUTPUT -p tcp --sport 25

${IPTABLES} -I INPUT  -p tcp --dport 465
${IPTABLES} -I OUTPUT -p tcp --sport 465

${IPTABLES} -I INPUT  -p tcp --dport 587
${IPTABLES} -I OUTPUT -p tcp --sport 587

${IPTABLES} -I INPUT  -p tcp --dport 995
${IPTABLES} -I OUTPUT -p tcp --sport 995

${IPTABLES} -I INPUT  -p tcp --dport 993
${IPTABLES} -I OUTPUT -p tcp --sport 993
