# configs/debian/systemd/bind9-resolvconf.service
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
Description=local BIND via resolvconf
Documentation=man:named(8) man:resolvconf(8)
Requires=bind9.service
After=bind9.service
ConditionFileIsExecutable=/sbin/resolvconf

[Service]
ExecStart=/bin/sh -c 'echo nameserver 127.0.0.1 | /sbin/resolvconf -a lo.named'
ExecStop=/sbin/resolvconf -d lo.named
RemainAfterExit=yes

[Install]
WantedBy=bind9.service
