# configs/debian/systemd/imscp_daemon.service
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
Description=i-MSCP communication daemon
After=local-fs.target network.target syslog.target

[Service]
Type=forking
PIDFile=/run/imscp_daemon.pid
ExecStart=/var/www/imscp/daemon/imscp_daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
