# configs/debian/systemd/imscp_traffic.service
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
Description=i-MSCP traffic logger
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/var/www/imscp/engine/bin/imscp-net-traffic add
ExecStop=/var/www/imscp/engine/bin/imscp-net-traffic remove

[Install]
WantedBy=multi-user.target
