# configs/Debian/systemd/imscp_daemon.service
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
Description=i-MSCP communication daemon
After=local-fs.target network.target syslog.target

[Service]
Type=forking
PIDFile=/var/run/imscp_daemon.pid
ExecStart=/var/www/imscp/daemon/imscp_daemon -b /var/www/imscp/engine/imscp-rqst-mngr -p /var/run/imscp_daemon.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
