# configs/debian/frontend/cron.daily/imscp_frontend
#!/bin/sh
# CRON(8) configuration file for i-MSCP frontEnd - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN

MAX_TMPFILE_LIFETIME=5
# MAX_TMPFILE_LIFETIME can be overridden in this file
if [ -r /etc/default/imscp_frontend ]; then
        . /etc/default/imscp_frontend
fi

if [ -d {WEB_DIR}/data/tmp ]; then
        find {WEB_DIR}/data/tmp -type f -mtime +$MAX_TMPFILE_LIFETIME -delete
fi
