# configs/debian/upstart/imscp_mountall.conf
# UPSTART(7) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
description "Mounts i-MSCP filesystems"
author      "Laurent Declercq <l.declercq@nuxwin.com>"

start on filesystem
stop on runlevel [016]

pre-start exec /usr/bin/perl /var/www/imscp/engine/tools/imscp-mountall start
post-stop exec /usr/bin/perl /var/www/imscp/engine/tools/imscp-mountall stop
