# configs/debian/init/imscp_daemon.conf
# imscp_daemon

description "i-MSCP communication daemon"
author      "Laurent Declercq <l.declercq@nuxwin.com>"

start on starting imscp_panel
stop on stopping imscp_panel

expect daemon
respawn

exec /var/www/imscp/daemon/imscp_daemon
