# configs/debian/sysvinit/imscp_mountall
#!/bin/sh
### BEGIN INIT INFO
# Provides:          imscp_mountall
# Required-Start:    $remote_fs
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Mounts i-MSCP filesystems
### END INIT INFO

# Note: do not modify any of these vars here, use /etc/default/$NAME instead

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
NAME=imscp_mountall
LFILE=/var/run/$NAME
TOOLSDIR="/var/www/imscp/engine/tools"

. /lib/lsb/init-functions

# Read config file if present.
if [ -r /etc/default/$NAME ]; then
    . /etc/default/$NAME
fi

check_for_upstart()
{
    # We don't use init_is_upstart() from /lib/lsb/init-functions because
    # that function is not provided in all supported Distro versions
    if [ -x /sbin/initctl ] && /sbin/initctl version 2>/dev/null | /bin/grep -q upstart; then
       exit $1
    fi
}

do_start()
{
    perl ${TOOLSDIR}/imscp-mountall start 2>/dev/null
    touch $LFILE
}

do_stop()
{
    perl ${TOOLSDIR}/imscp-mountall stop 2>/dev/null
    rm -f $LFILE
}


case "$1" in
  start)
    check_for_upstart 1
    log_daemon_msg "Starting $DESC" "$NAME"
    if [ -e "$LFILE" ]; then
        log_progress_msg "already started"
        log_end_msg 0
    else
      do_start
      log_end_msg $?
    fi
    ;;
  stop)
    check_for_upstart 0
    log_daemon_msg "Stopping $DESC" "$NAME"
    if [ ! -e "$LFILE" ]; then
        log_progress_msg "already stopped"
        log_end_msg 0
    else
        do_stop
        log_end_msg $?
    fi
    ;;
  restart|force-reload)
    check_for_upstart 1
    log_daemon_msg "Stopping $DESC" "$NAME"
    if [ ! -e "$LFILE" ]; then
        log_progress_msg "already stopped"
        log_end_msg 0
    else
        do_stop
        log_end_msg $?
    fi
    log_daemon_msg "Starting $DESC" "$NAME"
    do_start
    log_end_msg $?
    ;;
  status)
    check_for_upstart 1
    if [ -e "$LFILE" ]; then
        log_success_msg "$NAME is running"
    else
        log_failure_msg "$NAME is not running"
        exit 3
    fi
    ;;
  *)
    log_action_msg "Usage: /etc/init.d/$NAME {start|stop|restart|force-reload|status}" >&2
    exit 3
    ;;
esac

:
